# cygwin/etc/profile.d/xview.sh
export MANPATH="${MANPATH}:/usr/openwin/man"
export PATH="${PATH}:/usr/openwin/bin"
