# cygwin/etc/profile.d/xview.csh
if ( ! $?MANPATH ) setenv MANPATH ""
setenv MANPATH "${MANPATH}:/usr/openwin/man"
if ( ! $?PATH ) setenv PATH ""
setenv PATH "${PATH}:/usr/openwin/bin"
